# contrib/surface-pro-11-suspend/surface-resume-inhibit.service
[Unit]
Description=Inhibit suspend briefly after resume (Surface Pro 11)
After=systemd-suspend.service

[Service]
Type=simple
ExecStart=/usr/bin/systemd-inhibit --what=sleep --who="surface-resume" --why="Prevent immediate re-suspend after lid bounce" --mode=block /usr/bin/sleep 10

[Install]
WantedBy=suspend.target
