# contrib/surface-pro-11-suspend/surface-display-off.service
[Unit]
Description=Blank display before suspend (Surface Pro 11)
Before=systemd-suspend.service
After=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/surface-pre-suspend.sh

[Install]
WantedBy=sleep.target
