# root/etc/NetworkManager/dispatcher.d/no-power-save
#!/bin/sh

# MDJ
IFACE=$1
ACTION=$2

IW=/sbin/iw

test -x $IW || exit 0
[ "$IFACE" = "wlp2s0" ] || exit 0

case "$ACTION" in
    up)
	$IW dev $IFACE set power_save off
	;;
esac
