# pkg/fedora/kernel-surface/files/linux-surface-default-watchdog.service
[Unit]
Description=linux-surface default kernel watchdog

[Service]
Type=oneshot
ExecStart=/usr/bin/linux-surface-default-watchdog.py
