# pkg/fedora/kernel-surface/files/linux-surface-default-watchdog.path
[Unit]
Description=linux-surface default kernel watchdog

[Path]
PathChanged=/boot

[Install]
WantedBy=default.target
