# contrib/surface-pro-11-cameras/surface-camera-bridge@.service
[Unit]
Description=Surface Pro 11 Camera Bridge (%i)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/surface-camera-bridge.sh %i
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
