# pkg/arch/surface/surface-config.install
# arg 1:  the new package version
post_install() {
  echo "Surface Go users: Enable power management for Surface Go touchscreen"
  echo "    sudo systemctl enable -q surfacego-touchscreen"
}
